# .config/systemd/user/rotate-wallpapers.service
[Unit]
Description=Rotate wallpaper
Wants=rotate-wallpaper.timer

[Service]
ExecStart=/home/ezaurum/.fehbg

[Install]
WantedBy=multi-user.target
